# mnlp-django/jobs/crontab
0 * * * * cd ~/current/mnlp/ && python manage.py mailjob