# jobs/crontab
MAILTO=""
00 * * * * cd ~/current/mnlp/ && bash -lc "python manage.py dropboxjob"
05 * * * * cd ~/current/mnlp/ && bash -lc "python manage.py heatmapjob"
