# jobs/crontab
MAILTO=""
0 * * * * cd ~/current/mnlp/ && bash -lc "python manage.py dropboxjob"
