# jobs/crontab
MAILTO=""
30 * * * * cd ~/current/mnlp/ && bash -lc "python manage.py dropboxjob"
